# scripts/deploy/com.aguayluz.refresh-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  AguaYLuz daily refresh (USGS reservoir levels + export).
  Edit the two __ABSOLUTE_PATH__ placeholders, then:
    cp scripts/deploy/com.aguayluz.refresh-daily.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.aguayluz.refresh-daily.plist
  Runs every day at 06:15 local. Logs to logs/refresh-daily.{out,err}.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aguayluz.refresh-daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>__ABSOLUTE_PATH_TO_REPO__/.venv/bin/python</string>
    <string>__ABSOLUTE_PATH_TO_REPO__/scripts/refresh.py</string>
    <string>--daily</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__ABSOLUTE_PATH_TO_REPO__</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>EPA_WATERS_API_KEY</key>
    <string></string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>6</integer>
    <key>Minute</key><integer>15</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>__ABSOLUTE_PATH_TO_REPO__/logs/refresh-daily.out</string>
  <key>StandardErrorPath</key>
  <string>__ABSOLUTE_PATH_TO_REPO__/logs/refresh-daily.err</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
